# DEBIAN/postinst
#!/bin/bash
set -e

systemctl enable --now proton-ge-update.timer || true

exit 0

