# DEBIAN/postinst
#!/bin/bash
set -e

# Habilita o timer para todos usuários
systemctl --global enable proton-ge-update.timer || true
systemctl --global enable proton-ge-update.service || true

exit 0

