# usr/lib/systemd/user/proton-ge-update.service
[Unit]
Description=Atualizador automático do Proton GE

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-proton-ge.sh
