# usr/lib/systemd/system/proton-ge-update.timer
[Unit]
Description=Timer semanal para atualizar Proton GE

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

